# content/media/scripts/borg/borgstorage.timer
[Unit]
Description=Run borgstorage

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
